# certdash.service
[Unit]
Description=Certbot Certificate Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/letsencrypt
ExecStart=/usr/bin/python3 /etc/letsencrypt/certdash.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
